# etc/supervisord.d/rados-add-user.ini
[program:rados-add-user]
command=/opt/qnib/ceph/bin/create_rados_user.py --username "Test Qnib" qnib
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = false
startsecs=0
exitcodes=0
startretries=3
