# etc/supervisord.d/rados-add-user.ini
[program:rados-add-user]
command=/opt/qnib/ceph/bin/create_rados_user.py
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
